# install/shepard.service
[Unit]
Description=Shepard
After=login.target

[Service]
ExecStart=/home/uaarg/Documents/shepard/bin/run-nav.sh
WorkingDirectory=/home/uaarg/Documents/shepard

[Install]
WantedBy=multi-user.target
